# .config/systemd/user/session-cleanup.timer
[Unit]
Description=Session health maintenance timer
Requires=session-cleanup.service

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
WakeSystem=false
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
