# .config/systemd/user/shutdown-cleanup.service
[Unit]
Description=Clean shutdown for Hyprland session
DefaultDependencies=false
Before=shutdown.target reboot.target halt.target
After=hyprland-session.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStop=/bin/sh -c 'systemctl --user stop process-optimization.timer; pkill -TERM -f "process_optimization.sh"; pkill -TERM -f "systemd_integration.sh"'
TimeoutStopSec=10

[Install]
WantedBy=hyprland-session.target