# .config/systemd/user/waybar.service
#==========================================================#
#                      WAYBAR SERVICE                      #
#==========================================================#
# File: waybar.service
# Purpose: Highly customizable Wayland bar for Hyprland
# Last Updated: 2026-03-20
#==========================================================#

[Unit]
Description=Highly customizable Wayland bar for Hyprland
Documentation=man:waybar(5)
After=wayland-session@hyprland.desktop.target
PartOf=wayland-session@hyprland.desktop.target
Wants=pipewire.service

[Service]
Type=exec
ExecStart=/usr/bin/waybar
ExecReload=/bin/kill -SIGUSR2 $MAINPID
Restart=on-failure
RestartSec=3
KillMode=mixed
TimeoutStopSec=10
Nice=-2

[Install]
WantedBy=wayland-session@hyprland.desktop.target
